# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:bot]
user=root
command=bundle exec lita
directory=/app
autorestart=true
